# etc/inittab
# /etc/inittab - busybox init configuration

# system setup
::sysinit:/bin/ash /etc/rc

# login shell on serial console
ttyS0::respawn:-/bin/ash
